# thelounge/rootfs/etc/services.d/server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: The Lounge
# Runs The Lounge server
# ==============================================================================
export THELOUNGE_HOME=/data/thelounge

bashio::log.info "Starting The Lounge server.."

exec /usr/local/bin/thelounge start \
    -c port=15100 \
    -c https.enable="$(bashio::config 'ssl')" \
    -c https.ca="/ssl/$(bashio::config 'certfile')" \
    -c https.certificate="/ssl/$(bashio::config 'certfile')" \
    -c https.key="/ssl/$(bashio::config 'keyfile')" \
    -c fileUpload=true \
    -c prefetch=true \
    -c prefetchStorage=true \
    -c theme="$(bashio::config 'default_theme')"
