# thelounge/rootfs/etc/services.d/thelounge/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: The Lounge
# Runs The Lounge server
# ==============================================================================
export THELOUNGE_HOME=/data/thelounge

bashio::log.info "Starting The Lounge server.."

exec /usr/local/bin/thelounge start \
    -c fileUpload.enable=true \
    -c port=15100 \
    -c prefetch=true \
    -c prefetchStorage=true \
    -c reverseProxy=true \
    -c theme="$(bashio::config 'default_theme')"
