# thelounge/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: The Lounge
# Runs the Nginx daemon
# ==============================================================================

# Wait for thelounge to become available
bashio::net.wait_for 15100

bashio::log.info "Starting NGinx..."

exec nginx
