# thelounge/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: The Lounge
# Runs the Nginx daemon
# ==============================================================================
declare -a options

bashio::log.info "Starting Nginx..."

# Enable Hassio Ingress
export ENABLE_HASSIO_INGRESS=true

options+=(-g "daemon off;")

if bashio::config.true 'ssl'; then
    options+=(-c /etc/nginx/nginx-ssl.conf)
fi

exec nginx "${options[@]}"
